# scripts/rc.local
#!/bin/bash

# start habpi service
su habpi -c /home/habpi/payload/scripts/start &
